# roles/fs_snapshots/templates/zpool-scrub@.timer.j2
[Unit]
Description=Scrub ZFS pool %i

[Timer]
OnCalendar={{ zpool_scrub_schedule }}
Persistent=true
Unit=zpool-scrub@%i.service

[Install]
WantedBy=timers.target
