# roles/minecraft_fabric_server/templates/systemd/minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
User={{ minecraft_user }}
Nice=1
SuccessExitStatus=0 1
WorkingDirectory={{ minecraft_dir }}
ExecStart={{ minecraft_dir }}/start.sh
ExecStop={{ minecraft_dir }}/stop.sh

[Install]
WantedBy=multi-user.target
