# roles/minecraft_fabric_server/templates/systemd/minecraft-backup.service
[Unit]
Description=Minecraft Server backups
After=network.target

[Service]
User={{ minecraft_user }}
WorkingDirectory={{ minecraft_dir }}
ExecStart={{ minecraft_dir }}/backup.sh

[Install]
WantedBy=multi-user.target
